# wificheck.service
[Unit]
Description=Wifi Checking Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/DisneyStreaming
ExecStart=/home/pi/DisneyStreaming/wifi.sh

[Install]
WantedBy=multi-user.target